# systemd/liq-telegram-bot.service
[Unit]
Description=Liquidation Bot - Telegram Command Interface
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=ivan
WorkingDirectory=/home/ivan/liquidation-bot
ExecStart=/home/ivan/liquidation-bot/.venv/bin/python -m telegram_bot.app
Restart=always
RestartSec=30
EnvironmentFile=/home/ivan/liquidation-bot/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=liq-telegram-bot

[Install]
WantedBy=multi-user.target
